# service/wifi-kiosk.service
[Unit]
Description=WiFi Kiosk Backend Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi/wifi-kiosk/backend
ExecStart=/usr/bin/node /home/pi/wifi-kiosk/backend/server.js
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=wifi-kiosk-backend

# Zmienne środowiskowe
Environment="NODE_ENV=production"
Environment="PORT=80"

[Install]
WantedBy=multi-user.target
